# gh_pages/gnome-keyring-headless.service
[Unit]
Description=GNOME Keyring for Headless Operation (proton-drive-sync)
Before=proton-drive-sync.service

[Service]
Type=oneshot
ExecStart={{KEYRING_INIT_SCRIPT}}
RemainAfterExit=yes

[Install]
WantedBy=default.target
